# vita_monitor/supervisor_conf/conf.d/vita_monitor_hardware_sender.conf
; Embodied Intelligence 体征监控系统 — 硬件发送端 (Supervisor)
; 路径说明见 backend_receiver.conf 顶部。使用方式：同目录 backend_receiver.conf（mkdir logs、编辑 user/environment、supervisorctl reread/update）。
;
; 发送端常部署在具身设备上：若与后端/前端不在同一主机，请勿把 backend/frontend 的 conf 复制到该设备，
; 且不要复制 z-embodied-vita-monitor-group.conf（否则会引用不存在的 program）。
;
; 若程序 FATAL：查 supervisor_conf/logs/hardware_sender.log

[program:embodied-vita-monitor-hardware-sender]
directory=/ros2_ws/vita_monitor/hardware_sender
command=/bin/bash /ros2_ws/vita_monitor//supervisor_conf/scripts/start_hardware_sender.sh
user=root
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/ros2_ws/vita_monitor/supervisor_conf/logs/hardware_sender.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=EMBODIED_SENDER_PYTHON="/ros2_ws/vita_monitor/hardware_sender/.venv/bin/python"
